# packaging/app-proxy-report.cron
2 * * * * root /usr/sbin/proxy2db >/dev/null 2>&1
